# config/tesla/fleet_telemetry/systemd/tesla-telemetry-bridge.service
[Unit]
Description=Tesla Telemetry Bridge (tail feed → POST to Rails webhook)
After=fleet-telemetry.service
Requires=fleet-telemetry.service

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/apps/portfolio/current

ExecStart=/usr/bin/ruby bin/tesla_telemetry_bridge.rb

Environment=TELEMETRY_FEED=/var/log/tesla-telemetry/feed.jsonl
Environment=TELEMETRY_WEBHOOK=http://localhost:3141/webhooks/tesla_telemetry

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

NoNewPrivileges=true
ProtectSystem=strict
# Bridge reads the deployed script from /home/deploy/... so it needs
# /home visible. read-only is enough (we don't write there).
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
